# scripts/launchd/com.narimato.intelligence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.narimato.intelligence</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>REPLACE_WITH_REPO_PATH/scripts/guardian.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_REPO_PATH</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>REPLACE_WITH_REPO_PATH/logs/intelligence-guardian.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_REPO_PATH/logs/intelligence-guardian.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
